# 20181113_docker_cron_python/cron.d/crontest
# crontest
SHELL = /bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin
* * * * * root pwd >> /tmp/crontest.txt && python -c "import datetime;print(datetime.datetime.now())" >> /tmp/crontest.txt
